# config/sistema-ocr.service
[Unit]
Description=Sistema OCR FGJCDMX - Inicio Automático
Documentation=https://github.com/ELQUAyC/sistemaocr
After=network.target docker.service
Wants=docker.service
Requires=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
ExecStart=/home/vboxuser/Downloads/FJ1/scripts/inicio-sistema-inteligente.sh
ExecStop=/usr/bin/docker compose -f /home/vboxuser/Downloads/FJ1/docker-compose.yml down
WorkingDirectory=/home/vboxuser/Downloads/FJ1
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
Restart=no

# Variables de entorno
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=DOCKER_COMPOSE_FILE=/home/vboxuser/Downloads/FJ1/docker-compose.yml

[Install]
WantedBy=multi-user.target